# .config/systemd/user/backup-nutstore.timer
[Unit]
Description=Run backup-nutstore.service daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
