# .config/systemd/user/clean-cache.timer
[Unit]
Description=Run clean-cache.service weekly

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
