# .config/systemd/user/dotfiles-backup.timer
[Unit]
Description=Run dotfiles backup service hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
