# .config/systemd/user/downloads-sorter.service
[Unit]
Description=Downloads sorter service

[Service]
Type=oneshot
ExecStart=/usr/bin/true