# .config/systemd/user/file-sorter.service
[Unit]
Description=Sort files from download directories

[Service]
Type=oneshot
ExecStart=/home/peng/.config/systemd/user/file-sorter.sh
