# .config/systemd/user/rotate-wallpaper.service
[Unit]
Description=Change wallpaper periodically

[Service]
Type=oneshot
ExecStart=/home/peng/.config/niri/scripts/rotate-wallpaper-once.sh
